# supervisor/conf.d/start.conf
[program:tomcat]
command=bash /home/admin/bin/run.sh
user=root
environment=HOME="/home/admin",USER="admin"
